# profile.d/erlang.sh
# Output strings as unicode per default
# Enable shell history
export ERL_FLAGS="+pc unicode -kernel shell_history enabled"
