# profile.d/golang.sh
export GOPATH="$HOME/devel/gocode"
PATH="$GOPATH/bin:$PATH"
