# profile.d/golang.sh
export GOPATH="$HOME/development/gocode"
PATH="$GOPATH/bin:$PATH"
