# profile.d/less.sh
hash lesspipe && eval "$(lesspipe)"

export MANPAGER="less -FX"
