# profile.d/nodejs.sh
export NPM_CONFIG_PREFIX="$HOME/lib/nodejs"
export NODE_PATH="$NPM_CONFIG_PREFIX/lib/node_modules:/usr/lib/nodejs"
export PATH="$NPM_CONFIG_PREFIX/bin:$PATH"
