# profile.d/pager.sh
command -v lesspipe >/dev/null && eval "$(lesspipe)"

export MANPAGER="less -FX"
