# profile.d/pager.sh
command -v lesspipe >/dev/null && eval "$(lesspipe)"

export PAGER=less
export MANPAGER="less -FQX"
