# profile.d/plan9.sh
# Plan9port installation directory
PLAN9="$HOME/devel/plan9/plan9port"
export PLAN9
