# profile.d/qt.sh
qt_root=/opt/Qt
qt_version=5.8
qt_abi=gcc_64
export PATH="${qt_root}/${qt_version}/${qt_abi}/bin:${PATH}"
