# profile.d/ripgrep.sh
export RIPGREP_CONFIG_PATH="$HOME/lib/cfg/ripgreprc"
