# profile.d/rust.sh
export CARGO_HOME="$HOME/lib/cargo"
export RUSTUP_HOME="$HOME/lib/rustup"
# The bin path is exported by $CARGO_HOME/env
. "$CARGO_HOME/env"
