# profile.d/rust.sh
export CARGO_HOME="$HOME/lib/rust/cargo"
export RUSTUP_HOME="$HOME/lib/rust/rustup"
