# profile.d/stm32cube.sh
if test -d "$HOME/lib/apps/STM32CubeProgrammer/bin"
then export PATH="$HOME/lib/apps/STM32CubeProgrammer/bin:$PATH"
fi
