# profile.d/tmux_session.sh
export TMUX_SESSIONRC="$HOME/lib/cfg/tmux_sessionrc"
export TMUX_SESSION_HOME="$HOME/lib/cfg/tmux_session.d"
