# profile.d/tmux_session.sh
export TMUX_SESSIONRC="$HOME/lib/config/tmux_sessionrc"
export TMUX_SESSION_HOME="$HOME/lib/config/tmux_session.d"
